# service/runit/batts/run
#!/bin/sh

exec batts --interval-in-seconds=60 --log-directory=$XDG_DATA_HOME/batts
