# service/systemd/batstat-daemon.service
[Unit]
Description=batstat daemon

[Service]
Type=simple
ExecStart=/usr/bin/env batstat daemon

[Install]
WantedBy=default.target
